# deploy/systemd/user/shared-memory-monitor.service
# Shared Memory Monitor — systemd user unit (template)
# install-systemd-user.sh substitutes @MONITOR_ROOT@ with your checkout path.
# Installed to: ~/.config/systemd/user/shared-memory-monitor.service

[Unit]
Description=Shared Memory Monitor (poll loop + dashboard :8765)
Documentation=file://@MONITOR_ROOT@/README.md#running-as-a-service
After=network-online.target hive-mind-gateway.service
Wants=hive-mind-gateway.service

[Service]
Type=simple
WorkingDirectory=@MONITOR_ROOT@
EnvironmentFile=@MONITOR_ROOT@/.env
ExecStart=@MONITOR_ROOT@/.venv/bin/python -m sm_telemetry_monitor loop --serve --interval 600
Restart=on-failure
RestartSec=30

[Install]
WantedBy=default.target